# docker/supervisor/services/fix_permissions.conf
[program:fix_permissions]
command = /usr/sbin/fix_permissions.sh
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
exitcodes=0
startsecs=0