# docker/supervisor/services/generate_msmtprc.conf
[program:generate_msmtprc]
command = /usr/sbin/generate_msmtprc.sh
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
exitcodes=0
startsecs=0